# scripts/claw.service
[Unit]
Description=Claw Machine Job Queue Orchestrator
After=docker.service
Requires=docker.service

[Service]
Type=simple
User=claw
Group=claw
WorkingDirectory=/opt/clawmachine
EnvironmentFile=/opt/clawmachine/.env
ExecStart=/usr/bin/docker compose -f /opt/clawmachine/docker/docker-compose.yml --env-file /opt/clawmachine/.env up
ExecStop=/usr/bin/docker compose -f /opt/clawmachine/docker/docker-compose.yml --env-file /opt/clawmachine/.env down
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
